# src/99-pwg_initdb.conf
[group:piwigo]
programs=pwgdb

[program:pwgdb]
command = /usr/bin/curl -s -f -X POST --data-raw language=en_GB&dbhost=%(ENV_PWG_DB_HOST)s&dbuser=%(ENV_PWG_DB_USER)s&dbpasswd=%(ENV_PWG_DB_PASSWORD)s&dbname=%(ENV_PWG_DB_BASE)s&prefix=%(ENV_PWG_DB_PREFIX)s&admin_name=%(ENV_PWG_ADMIN_USERNAME)s&admin_pass1=%(ENV_PWG_ADMIN_PASSWORD)s&admin_pass2=%(ENV_PWG_ADMIN_PASSWORD)s&admin_mail=%(ENV_PWG_ADMIN_EMAIL)s&install=Start+installation --compressed --insecure http://localhost/piwigo/install.php?language=en_GB
process_name=%(program_name)s
priority=4
startsecs=0
autostart = true
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


